# meshdesk/MESHdesk/md_prerun
#!/bin/sh /etc/rc.common
# ============================================
# == Coova Chilli Startup Script =============
# == Hardware: Ubiquity PicoStation2 =========
# == OpenWRT version 12.09 ===================
# == Date 2013-07-03 =========================
# ============================================
 
START=50
STOP=50
DIR='/etc/MESHdesk'

start() {

	#Due to problems with the package postinstall script we will do a first run here
	FIRST_RUN=`uci get meshdesk.settings.first_run`
	echo "$FIRST_RUN"
	if [ $FIRST_RUN == 1 ]; then
		echo "First run - do some housekeeping"
		if [ -e /etc/init.d/alfred ]; then
                    echo "Disable rc.d symlink for Alfred"
        	    /etc/init.d/alfred disable
		    /etc/init.d/alfred stop
		   cp /etc/MESHdesk/files/common/alfred /etc/init.d/alfred
		fi

		if [ -e /etc/init.d/odhcpd ]; then
		    echo "Disable rc.d symlink for Odhcpd"
		    /etc/init.d/odhcpd disable
		    /etc/init.d/odhcpd stop
		fi
		
		if [ -e /etc/init.d/dnsmasq ]; then
		    echo "Disable rc.d symlink for DNSMASQ"
		    /etc/init.d/dnsmasq disable
		    /etc/init.d/dnsmasq stop
		fi
		
		if [ -e /etc/init.d/chilli ]; then
		    echo "Disable rc.d symlink for CoovaChilli"
		    /etc/init.d/chilli disable
		    /etc/init.d/chilli stop
		fi
		
		if [ -e /etc/init.d/privoxy ]; then
		    echo "Disable rc.d symlink for Privoxy"
		    /etc/init.d/privoxy disable
		    /etc/init.d/privoxy stop
		fi
		#FIXME REMOVE THIS LATER
		if [ -d /lib/firmware/ath10k/QCA988X/hw2.0 ]; then
		    echo "Copy fixes ower for Archer board"
		    cp /etc/MESHdesk/files/common/802_11_ac/mac80211.sh /lib/netifd/wireless/mac80211.sh
		    cp /etc/MESHdesk/files/common/802_11_ac/firmware-5.bin /lib/firmware/ath10k/QCA988X/firmware-5.bin
		    cp /etc/MESHdesk/files/common/802_11_ac/55-ath10k /etc/modules.d/55-ath10k
		fi
		
		uci set meshdesk.settings.first_run=0
		uci commit
	else
		echo "First run completed..."
	fi
	
    	echo start
    	cd $DIR	
    	./a.lua &
}

stop() {
    echo stop
    killall /etc/MESHdesk/a.lua
}
